# scripts/magento/supervisor/magento-consumers.conf
[group:magento_consumers]
programs=saveconfig,product_action,product_action_website,catalog_website_sync,media_resize,export,inventory_cleanup,inventory_update,inventory_reservations_cleanup,inventory_reservations_update,inventory_salability,inventory_sourceitem,inventory_stock,media_content,media_gallery_renditions,media_gallery_sync,codegenerator,sales_rule_coupon,sales_rule_quote,product_alert,async_operations

[program:saveconfig]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start saveConfigProcessor --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/saveconfig.log

[program:product_action]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start product_action_attribute.update --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/product_action.log

[program:product_action_website]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start product_action_attribute.website.update --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/product_action_website.log

[program:catalog_website_sync]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start catalog_website_attribute_value_sync --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/catalog_website_sync.log

[program:media_resize]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start media.storage.catalog.image.resize --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/media_resize.log

[program:export]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start exportProcessor --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/export.log

[program:inventory_cleanup]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.source.items.cleanup --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_cleanup.log

[program:inventory_update]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.mass.update --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_update.log

[program:inventory_reservations_cleanup]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.reservations.cleanup --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_reservations_cleanup.log

[program:inventory_reservations_update]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.reservations.update --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_reservations_update.log

[program:inventory_salability]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.reservations.updateSalabilityStatus --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_salability.log

[program:inventory_sourceitem]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.indexer.sourceItem --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_sourceitem.log

[program:inventory_stock]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start inventory.indexer.stock --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/inventory_stock.log

[program:media_content]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start media.content.synchronization --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/media_content.log

[program:media_gallery_renditions]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start media.gallery.renditions.update --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/media_gallery_renditions.log

[program:media_gallery_sync]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start media.gallery.synchronization --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/media_gallery_sync.log

[program:codegenerator]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start codegeneratorProcessor --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/codegenerator.log

[program:sales_rule_coupon]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start sales.rule.update.coupon.usage --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/sales_rule_coupon.log

[program:sales_rule_quote]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start sales.rule.quote.trigger.recollect --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/sales_rule_quote.log

[program:product_alert]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start product_alert --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/product_alert.log

[program:async_operations]
command=/usr/bin/php /var/www/html/bin/magento queue:consumers:start async.operations.all --single-thread
autostart=true
autorestart=true
user=doge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/magento/consumers/async_operations.log 